# deploy/systemd/repo2-rag-jobs.service
[Unit]
Description=Repo2 RAG Jobs
[Service]
Type=simple
ExecStart=/usr/bin/bash scripts/rebuild-index.sh
